# meta-kudo/recipes-phosphor/ipmi/phosphor-ipmi-ssif/ssifbridge.service
[Unit]
Description=SSIF Bridge Daemon

[Service]
Type=dbus
BusName=org.openbmc.HostIpmi
ExecStart=/usr/sbin/ssifbridged
Restart=always

[Install]
WantedBy=multi-user.target
